# .dev/docker/root/etc/cron.d/pterodactyl-scheduler
* * * * * php /app/artisan schedule:run >> /dev/null 2>&1

# renew letsencrypt certs
0 0,12 * * * python -c 'import random; import time; time.sleep(random.random() * 3600)' && certbot renew 
